# vps-monitor.service
[Unit]
Description=VPS 库存/补货监控 (Telegram 推送)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vps-monitor
ExecStart=/usr/bin/python3 /opt/vps-monitor/monitor.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/vps-monitor/.env

# 日志: 直接走 python 的 FileHandler, journalctl 留作辅助
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
